# i3_dots/autostart/feh.desktop
[Desktop Entry]
Exec=sh -c $HOME/.fehbg
Name=feh
Type=Application
